# ethereum/node/supervisord.conf
[supervisord]
logfile = /var/log/supervisor/supervisord.log
loglevel = DEBUG

[program:geth]
command=sh /data/start_geth.sh
redirect_stderr=true
priority=2
